# files/etc/init.d/f50-wan-switch
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
	/root/f50-wan-switch.sh auto >/dev/null 2>&1 || true
}
